# packer/files/myunit.service
[Unit]
Description=My Service
 
 
[Service]
Type=simple
WorkingDirectory=/home/filip_filipov/reddit
ExecStart=/bin/bash -lc 'puma'
Restart=always
